# systemd/codex-waybar.service
[Unit]
Description=Codex Waybar Cache Writer
Documentation=https://github.com/jtaw5649/codex-waybar
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/codex-waybar --cache-file %h/.cache/codex-waybar/latest.json --max-chars 110 --session-window 6 --poll-ms 100
Restart=on-failure
RestartSec=0.3
StartLimitIntervalSec=10
StartLimitBurst=10
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
